# packages/vps/unified/supervisord.conf
[supervisord]
nodaemon=true
logfile=/data/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
pidfile=/var/run/supervisord.pid
user=root

[program:postgres]
command=/usr/local/bin/pg-start
user=postgres
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/data/logs/postgres.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/data/logs/postgres-error.log
stderr_logfile_maxbytes=5MB
priority=10

[program:redis]
command=/usr/bin/redis-server /etc/redis/redis.conf
user=redis
autostart=true
autorestart=true
startretries=3
startsecs=2
stdout_logfile=/data/logs/redis.log
stdout_logfile_maxbytes=2MB
stderr_logfile=/data/logs/redis-error.log
stderr_logfile_maxbytes=2MB
priority=20

[program:moltclip]
command=/usr/local/bin/moltclip-start
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=10
stdout_logfile=/data/logs/moltclip.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_logfile=/data/logs/moltclip-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
priority=30
